# etc/stream_list.timer
[Unit]
Description=Execute stream_list every 5 mins

[Timer]
OnCalendar=*:0/5
Unit=stream_list.service

[Install]
WantedBy=multi-user.target
